# debian/circus.upstart
description "circusd"

start on filesystem and net-device-up IFACE=lo
stop on runlevel [016]

respawn
exec /usr/bin/circusd /etc/circus/circusd.ini
